# systemd/daily-discover-flask.service.example
[Unit]
Description=Daily Discover Flask Backend
After=network.target postgresql.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/project
ExecStart=/path/to/uv run gunicorn asgi:application \
    -w 2 \
    -k uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8080 \
    --access-logfile /var/log/daily-discover/flask-access.log \
    --error-logfile /var/log/daily-discover/flask-error.log

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
