# systemd/daily-discover-flask.service
[Unit]
Description=Daily Discover Flask Backend
After=network.target postgresql.service

[Service]
Type=simple
User=gremlin
WorkingDirectory=/home/gremlin/blagh
ExecStart=/home/gremlin/.local/bin/uv run gunicorn app:app \
    -w 2 \
    --bind 0.0.0.0:8080 \
    --access-logfile - \
    --error-logfile -

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
